# scripts/systemd/llama-skein-update.service
# llama-skein unattended self-update (source rebuild + health-checked swap).
# Install: copy to /etc/systemd/system/ (system scope) and edit the Environment
# lines to match the host (binary path, service name, scope, health URL).
# See scripts/self-update.sh for all env knobs.
#
# Pair with llama-skein-update.timer. Enable with:
#   sudo systemctl daemon-reload && sudo systemctl enable --now llama-skein-update.timer
[Unit]
Description=llama-skein source self-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Host-specific overrides. Known fleet profiles (verify before enabling):
#
#   proxmox LXC 1016 — system service "llama-swap", listens :8080, binary
#   /usr/local/bin/llama-swap -> llama-skein. NOTE: go must be on PATH for the
#   self-build (it was not on the container's non-login PATH at last check — set
#   PATH or point LLAMA_SKEIN at /usr/local/go/bin). Use:
#     Environment=LLAMA_SKEIN_SERVICE=llama-swap
#     Environment=LLAMA_SKEIN_BIN=/usr/local/bin/llama-swap
#     Environment=LLAMA_SKEIN_HEALTH_URL=http://127.0.0.1:8080/v1/models
#     Environment=PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
#
#   rocky — has NO Go toolchain; cannot self-build. Drive its update from a Mac
#   build host (cross-compile + scp), not this unit.
Environment=LLAMA_SKEIN_BIN=/usr/local/bin/llama-skein
Environment=LLAMA_SKEIN_SERVICE=llama-skein
Environment=LLAMA_SKEIN_SERVICE_SCOPE=system
Environment=LLAMA_SKEIN_HEALTH_URL=http://127.0.0.1:11435/v1/models
ExecStart=/usr/local/bin/llama-skein-self-update.sh
